# not.conf
description "Run app"

setuid deploy
setgid deploy

start on runlevel startup
stop on runlevel shutdown

respawn

chdir /var/notary

script
        exec java -jar /var/notary/notary.jar
end script  
